# build/centos_x86_64/redsocks.service
[Unit]
Description=Transparent redirector of any TCP connection to proxy using your firewall

[Service]
Type=forking
PIDFile=/run/redsocks/redsocks.pid
EnvironmentFile=/etc/
User=redsocks
ExecStartPre=/usr/bin/redsocks -t
ExecStart=/usr/bin/redsocks -p /run/redsocks/redsocks.pid
ExecStopPost=/bin/rm /run/redsocks/redsocks.pid
Restart=on-abort

[Install]
WantedBy=multi-user.target
