# build/alpine_x86_64/install.sh
cp ./redsocks /usr/bin/redsocks
cp ./redsocks.conf /etc/redsocks.conf
cp ./redsocksd /etc/init.d/redsocksd
rc-update add redsocks2 boot